# crontab
# 6:00 AM todos los días
0 6 * * * echo "Cron ejecutado a las 06:00 $(date)" && curl -s http://localhost:8080/publicar_ahora

# 6:00 PM todos los días
0 18 * * * echo "Cron ejecutado a las 18:00 $(date)" && curl -s http://localhost:8080/publicar_ahora

