# crontab
00 18 * * * curl -s https://pos-feed-instagram.fly.dev/publicar_ahora
