# crontab
* * * * * echo "Hola soy cron $(date)" >> /proc/1/fd/1
